# contrib/loupe-worker.service
[Unit]
Description=loupe scan and verification worker
Documentation=https://github.com/tnull/loupe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=loupe-worker
Group=loupe-worker
EnvironmentFile=-/etc/loupe/loupe-worker.env
WorkingDirectory=/var/lib/loupe-worker
ExecStart=/opt/loupe/bin/loupe-worker run
Restart=always
RestartSec=5s
KillSignal=SIGINT
TimeoutStopSec=60s

# The worker invokes git plus LLM CLIs and may use bubblewrap's setuid
# helper on distributions without unprivileged user namespaces. Avoid
# NoNewPrivileges/ProtectHome here; keep the unit compatible with those
# tools and put credentials under the service user's own home or env file.
StateDirectory=loupe-worker
CacheDirectory=loupe-worker
RuntimeDirectory=loupe-worker
UMask=0077
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/loupe-worker /var/cache/loupe-worker

[Install]
WantedBy=multi-user.target
