# contrib/loupe-server.service
[Unit]
Description=loupe security scanning server
Documentation=https://github.com/tnull/loupe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=loupe
Group=loupe
EnvironmentFile=-/etc/loupe/loupe-server.env
# Secrets such as LOUPE_MASTER_KEY and TLS PEMs are injected by deploy-server.sh
# through systemd's runtime manager environment for each explicit
# restart, rather than persisted in /etc/loupe/loupe-server.env.
PassEnvironment=LOUPE_MASTER_KEY LOUPE_SERVER_CERT_PEM LOUPE_SERVER_CERT_PEM_B64 LOUPE_SERVER_KEY_PEM LOUPE_SERVER_KEY_PEM_B64 LOUPE_CA_CERT_PEM LOUPE_CA_CERT_PEM_B64 LOUPE_CA_KEY_PEM LOUPE_CA_KEY_PEM_B64
WorkingDirectory=/var/lib/loupe
ExecStart=/opt/loupe/bin/loupe-server serve
Restart=on-failure
RestartSec=5s
KillSignal=SIGINT
TimeoutStopSec=30s

# systemd creates /var/lib/loupe when the service starts. Keep secrets
# unreadable by other local users; root reads EnvironmentFile before
# dropping to User=.
StateDirectory=loupe
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/loupe

[Install]
WantedBy=multi-user.target
